# deploy/supervisor/wna-production.conf
[program:wna-production]
command=/opt/wna/env/bin/python /opt/wna/wna/service/web.py
command=/opt/wna/env/bin/gunicorn -c /opt/wna/wna/deploy/gunicorn_config.py service.web:app
environment=APP_CONFIG=/opt/wna/wna/config/rootcfg.py
user=cloo
directory=/opt/wna/wna/
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
